# dpkg-server/DEBIAN/prerm
#!/bin/sh
# prerm - removal script for minecontrold package

update-rc.d -f minecontrold remove
