# minecontrol-deb/DEBIAN/postinst
#!/bin/sh
# postinst - post install script for minecontrold package

update-rc.d -f minecontrold defaults
